# deploy/facestack.service
# FaceStack systemd unit (user service — no root needed).
#
# Install on motis:
#   mkdir -p ~/.config/systemd/user
#   cp deploy/facestack.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now facestack
#   # so it runs without an active login / starts at boot (one-time, needs sudo):
#   sudo loginctl enable-linger aras
#
# Manage:  systemctl --user {status,restart,stop} facestack
# Logs:    journalctl --user -u facestack -f

[Unit]
Description=FaceStack face recognition service
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=%h/FaceStack
# ROCm compat shim so ROCMExecutionProvider loads (see scripts/setup_rocm_motis.sh)
Environment=LD_LIBRARY_PATH=%h/rocm-compat
ExecStart=%h/FaceStack/.venv/bin/uvicorn facestack.service.app:app --host 0.0.0.0 --port 8011
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
